# pkgroot/Library/LaunchDaemons/com.clburlison.pinpoint.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clburlison.pinpoint</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Application Support/pinpoint/bin/supervisor</string>
    <string>--delayrandom</string>
    <string>300</string>
    <string>--timeout</string>
    <string>300</string>
    <string>--</string>
    <string>/Library/Application Support/pinpoint/bin/pinpoint</string>
    <string>--auto</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
</dict>
</plist>
